# etc/systemd/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin aidan --skip-login --nonewline --noissue --noreset --noclear - ${TERM}
