# etc/systemd/getty@tty1.service.d/autologin.conf
# YOUR_USER is substituted at apply time — see etc/README.md, don't plain-cp this
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin YOUR_USER --skip-login --nonewline --noissue --noreset --noclear - ${TERM}
